# templates/savanna-api.conf.erb
description "Savanna API Server upstart job"
author "Andy Edmonds <andy@edmonds.be>"

start on runlevel [2345]
stop on runlevel [!2345]


chdir <%= scope.lookupvar("savanna::params::sys_rundir") %>

pre-start script
    mkdir -p <%= scope.lookupvar("savanna::params::savanna_rundir") %>
    chown savanna:root <%= scope.lookupvar("savanna::params::savanna_rundir") %>

    mkdir -p <%= scope.lookupvar("savanna::params::savanna_lockdir") %>
    chown savanna:root <%= scope.lookupvar("savanna::params::savanna_lockdir") %>
end script

exec start-stop-daemon --start --chuid savanna --exec /usr/local/bin/<%= scope.lookupvar("savanna::params::savanna_service") %> -- --config-file=<%= scope.lookupvar("savanna::params::savanna_conf_file") %>